# examples/example.crontab
# This example cron job will run at 10 AM every Monday
# For more information on cron syntax, you can use: https://crontab-generator.org/
0 10 * * 1 root /usr/bin/python3 /root/mattermost-coffeebot/pair.py >> /var/log/coffeebot.log 2>&1
